# s6/jmx_exporter/run
#!/usr/bin/execlineb -P
java
-XX:+UseContainerSupport
-XshowSettings:vm
-jar /usr/local/jmx_prometheus/jmx_prometheus_httpserver.jar
5556
/usr/local/jmx_prometheus/jmx-kafka-connect-prometheus.yml
