# s6/nginx/finish
#!/usr/bin/execlineb -S0

foreground {
  echo "nginx finished, stopping everything"
}
s6-svscanctl -t /var/run/s6/services
