# s6/connect/finish
#!/usr/bin/execlineb -S0

foreground {
  echo "connect finished, stopping everything"
}
s6-svscanctl -t /var/run/s6/services
